# build/config-leader/runit/app/run
#!/bin/sh
DEBUG="*pryv*" exec chpst -u app node /app/bin/src/server.js --config /app/conf/config-leader.json
